# flashh-stays.timer
[Unit]
Description=Flashh-stays timer

[Timer]
OnCalendar=*-*-* 8..23:0/15
#OnCalendar=*:0/1

[Install]
WantedBy=multi-user.target
